# package/base-files/files/etc/uci-defaults/system_defaults
#!/bin/sh

. /lib/functions/uci-defaults.sh
. /lib/platform.sh

tw_set_dhcpc

uci commit system

exit 0
