# target/linux/ralink/mt7620/initramfs/common/files/etc/inittab
::sysinit:/etc/rcS S boot
::shutdown:/etc/rcS K stop
ttyS1::askfirst:/bin/ash --login
