# staging_dir/target-mipsel_r2_uClibc-0.9.33.2/root-ralink/etc/hotplug.d/wifi/10-normol
#!/bin/sh

echo "$ACTION" "$VAP" "$MACADDR" "$RSSI" > /tmp/wifi.txt

