# staging_dir/target-mipsel_r2_uClibc-0.9.33.2/root-ralink/etc/init.d/hcsh
#!/bin/sh /etc/rc.common

START=10
STOP=90

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

start() {
	stop
	service_start /usr/sbin/hcshd
}

stop() {
	service_stop /usr/sbin/hcshd
}

