# staging_dir/target-mipsel_r2_uClibc-0.9.33.2/root-ralink/etc/init.d/leds-stats
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2013 hiwifi.com
#

START=99

start() {

	if [ -e /etc/config/led_disable ]; then
		echo 0 > /proc/hiwifi/eth_led
		setled off green system
		setled off green internet
	fi	
}
