# staging_dir/target-mipsel_r2_uClibc-0.9.33.2/root-ralink/etc/init.d/wisp_chk
#!/bin/sh /etc/rc.common

START=99

start() {
	wisp_chk.sh &
}

stop() {
	killall wisp_chk.sh
}
