# staging_dir/target-mipsel_r2_uClibc-0.9.33.2/root-ralink/etc/hotplug.d/gmac/wan/00-internet-led
#!/bin/sh /etc/rc.common

stop() {
	wan_ifname=$(uci get network.wan.ifname)

	if [ $1 = $wan_ifname ]; then
		setled off green internet
	fi
}
