# staging_dir/target-mipsel_r2_uClibc-0.9.33.2/root-ralink/etc/hotplug.d/gmac/wan/02-dhcp-operation
#!/bin/sh /etc/rc.common

stop()
{
	if [ "dhcp" = "$(uci get network.wan.proto)" ];then
		kill -USR2 $(pidof udhcpc)
	fi
}

start()
{
	if [ "dhcp" = "$(uci get network.wan.proto)" ];then
	        kill -USR1 $(pidof udhcpc)
	fi
}
